# dotfiles/systemd/user/duck.service
[Unit]
Description=Update dyndns

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/home/$USER/env/bin/duckdns'
