# scripts/systemd-services/duck.service
[Unit]
Description=Update dyndns

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/local/bin/duck.sh'
